# services/client-updates
#!/sbin/runscript
# Copyright 2009 Fabio Erculiani
# Distributed under the terms of the GNU General Public License v2

DAEMON_EXEC="/usr/sbin/client-updates"
PID="/var/run/entropy_client_daemon.pid"
CMDLINE=""

depend() {
    use dbus
}

start() {

    ebegin "Starting Entropy Client daemon"
    start-stop-daemon --background --make-pidfile --start --pidfile ${PID} --quiet --exec ${DAEMON} -- ${CMDLINE}
    sleep 1
    if [ -f "${PID}" ]; then
        if [ -d "/proc/$(cat ${PID})" ]; then
            eend 0
        else
            eend 1
        fi
    else
        eend 1
    fi

}

stop() {

    ebegin "Stopping Entropy Client daemon"
    if [ -f "${PID}" ]; then
        mypid=$(cat $PID 2&> /dev/null)
        if [ -d "/proc/$(cat ${mypid})" ]; then
            eend 0
        else
            eend 1
        fi
    else
        eend 1
    fi

}
